# systemd/blockhost-admin.service
[Unit]
Description=BlockHost Admin Panel
After=network-online.target blockhost-root-agent.service
Wants=network-online.target

[Service]
Type=simple
User=blockhost
Group=blockhost
WorkingDirectory=/opt/blockhost
ExecStart=/usr/bin/python3 -m admin.app --host 0.0.0.0 --port 8443
Restart=on-failure
RestartSec=5
EnvironmentFile=-/opt/blockhost/.env

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
