# systemd/blockhost-firstboot.service
[Unit]
Description=BlockHost First-Boot Installer
After=network-online.target
Wants=network-online.target
ConditionPathExists=!/var/lib/blockhost/.setup-complete

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/opt/blockhost/first-boot.sh
TimeoutStartSec=3600
KillMode=process

# Output to journal and let script handle TTY
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
